# infrastructure/terraform/templates/crucible-updater.timer
[Unit]
Description=Check for Crucible Platform updates every 5 minutes
Requires=crucible-updater.service

[Timer]
# Run 2 minutes after boot
OnBootSec=2min
# Then every 5 minutes
OnUnitActiveSec=5min
# Randomize by up to 30 seconds to avoid thundering herd
RandomizedDelaySec=30

[Install]
WantedBy=timers.target